# krbsnap/krbsnap.crontab
# m h  dom mon dow   command
*/5 *  *   *   *     ./sync_gitsvn >> sync.log 2>&1
  0 3  *   *   *     (cd snap && ./gensnap) > /dev/null 2>&1
